# PostgreSQL16/rootfs/etc/services.d/exporter/run
#!/usr/bin/with-contenv bashio

DB_USER=$(bashio::config 'postgres_user')
DB_PASS=$(bashio::config 'postgres_password')
# Connect to localhost since we are in the same container
DATA_SOURCE_NAME="postgresql://${DB_USER}:${DB_PASS}@localhost:5432/postgres?sslmode=disable"
export DATA_SOURCE_NAME

bashio::log.info "Starting Prometheus Postgres Exporter..."
# Wait for postgres to be ready? Exporter usually retries.
# optional: wait-for-it equivalent
until nc -z localhost 5432; do
    bashio::log.debug "Waiting for PostgreSQL..."
    sleep 1
done

exec prometheus-postgres-exporter \
    --web.listen-address=":9187"
