# PostgreSQL18/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio

# Set up data directory
export PGDATA=/data/postgresql
mkdir -p "$PGDATA"
chown -R postgres:postgres "$PGDATA"
chmod 0700 "$PGDATA"

# Set up runtime directory
mkdir -p /run/postgresql
chown -R postgres:postgres /run/postgresql

# Initialize database if it doesn't exist
if ! bashio::fs.file_exists "$PGDATA/PG_VERSION"; then
    bashio::log.info "Initializing PostgreSQL database..."
    su-exec postgres initdb -D "$PGDATA"
    
    # Configure access - allow all from local network (S6 container networking)
    echo "host all all 0.0.0.0/0 md5" >> "$PGDATA/pg_hba.conf"
    
    # Listen on all interfaces
    echo "listen_addresses = '*'" >> "$PGDATA/postgresql.conf"
fi

# Get config
DB_USER=$(bashio::config 'postgres_user')
DB_PASS=$(bashio::config 'postgres_password')
CONFIG_CONTENT=$(bashio::config 'postgres_config')

# Apply custom configuration before starting
if [ -n "$CONFIG_CONTENT" ]; then
    bashio::log.info "Applying custom PostgreSQL configuration."
    AUTO_CONF="$PGDATA/postgresql.auto.conf"
    echo "$CONFIG_CONTENT" > "$AUTO_CONF"
    chown postgres:postgres "$AUTO_CONF"
fi

# Start PostgreSQL
bashio::log.info "Starting PostgreSQL..."
su-exec postgres postgres -D "$PGDATA" &
PG_PID=$!

# Wait for PostgreSQL to be ready
until su-exec postgres pg_isready; do
    sleep 1
done

# Create user if it doesn't exist
if ! su-exec postgres psql -tAc "SELECT 1 FROM pg_roles WHERE rolname='${DB_USER}'" | grep -q 1; then
    bashio::log.info "Creating user ${DB_USER}..."
    su-exec postgres psql -c "CREATE USER \"${DB_USER}\" WITH PASSWORD '${DB_PASS}';"
else
    bashio::log.info "User ${DB_USER} already exists."
fi

# Create/verify databases
bashio::log.info "Checking databases..."
for db in $(bashio::config 'postgres_db'); do
    if su-exec postgres psql -lqt | cut -d \| -f 1 | grep -qw "$db"; then
        bashio::log.info "Database ${db} exists."
    else
        bashio::log.info "Creating database ${db}..."
        su-exec postgres psql -c "CREATE DATABASE \"${db}\" OWNER \"${DB_USER}\";"
    fi
done

bashio::log.info "PostgreSQL is ready!"

# Wait for PostgreSQL process
wait "$PG_PID"
