# n8n/rootfs/etc/services.d/n8n/run
#!/usr/bin/with-contenv bashio

# Set up data directory
export N8N_USER_FOLDER="/data/n8n"
mkdir -p "$N8N_USER_FOLDER"

# Network settings
export N8N_PORT="5678"
export N8N_HOST="0.0.0.0"
export N8N_SECURE_COOKIE="$(bashio::config 'n8n_secure_cookie')"

# Ingress configuration
INGRESS_URL="$(bashio::addon.ingress_url)"
export N8N_EDITOR_BASE_URL="${INGRESS_URL}"
export WEBHOOK_URL="${INGRESS_URL}"
export N8N_PATH="${INGRESS_URL#*://*/}"  # Extract path from ingress URL

# Disable telemetry
export N8N_DIAGNOSTICS_ENABLED=false
export N8N_PERSONALIZATION_ENABLED=false
export N8N_ENFORCE_SETTINGS_FILE_PERMISSIONS=false
export N8N_RUNNERS_ENABLED=true
export N8N_BLOCK_ENV_ACCESS_IN_NODE=false
export N8N_DISABLE_ENV_ACCESS_IN_NODE=true
export N8N_GIT_NODE_DISABLE_BARE_REPOS=true

# Database configuration
DB_TYPE=$(bashio::config 'db_type')
export DB_TYPE="${DB_TYPE}"

if [ "${DB_TYPE}" = "postgresdb" ]; then
    bashio::log.info "Configuring PostgreSQL database..."
    export DB_POSTGRESDB_HOST="$(bashio::config 'db_postgresdb_host')"
    export DB_POSTGRESDB_PORT="$(bashio::config 'db_postgresdb_port')"
    export DB_POSTGRESDB_DATABASE="$(bashio::config 'db_postgresdb_database')"
    export DB_POSTGRESDB_USER="$(bashio::config 'db_postgresdb_user')"
    export DB_POSTGRESDB_PASSWORD="$(bashio::config 'db_postgresdb_password')"
else
    bashio::log.info "Using SQLite database..."
    export DB_SQLITE_POOL_SIZE=10
fi


# Authentication
if bashio::config.true 'n8n_basic_auth_active'; then
    export N8N_BASIC_AUTH_ACTIVE="true"
    export N8N_BASIC_AUTH_USER="$(bashio::config 'n8n_basic_auth_user')"
    export N8N_BASIC_AUTH_PASSWORD="$(bashio::config 'n8n_basic_auth_password')"
else
    export N8N_BASIC_AUTH_ACTIVE="false"
fi

bashio::log.info "Starting n8n..."
exec n8n start
